# daily_notes/src/crontab
*/15 * * * * python3 /src/sleep_cal.py
10 */6 * * * python3 /src/submit_day.py
5-59/15 * * * * python3 /src/watchtime.py
